# ansible/roles/monitoring_agents/templates/secos-agents.service.j2
[Unit]
Description=SecOS Monitoring Agents (Promtail + Node Exporter, rootless)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
WorkingDirectory=%h/secos/agents
EnvironmentFile=%h/secos/agents/.env
ExecStartPre=-/usr/bin/podman-compose -f %h/secos/agents/compose.yml down
ExecStart=/usr/bin/podman-compose -f %h/secos/agents/compose.yml up
ExecStop=/usr/bin/podman-compose -f %h/secos/agents/compose.yml down
TimeoutStartSec=120
TimeoutStopSec=60

[Install]
WantedBy=default.target
