# ansible/roles/monitoring_stack/files/secos-monitoring.service
[Unit]
Description=SecOS Monitoring Stack (Prometheus + Loki + Grafana, rootless)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Restart=on-failure
RestartSec=10
WorkingDirectory=%h/secos/monitoring
ExecStartPre=-/usr/bin/podman-compose -f %h/secos/monitoring/compose.yml down
ExecStart=/usr/bin/podman-compose -f %h/secos/monitoring/compose.yml up
ExecStop=/usr/bin/podman-compose -f %h/secos/monitoring/compose.yml down
TimeoutStartSec=300
TimeoutStopSec=60

[Install]
WantedBy=default.target
